# templates/vector.service.j2
# Vector systemd service file
# Managed by Ansible
# Generated from template: vector.service.j2

[Unit]
Description=Vector
Documentation=https://vector.dev
After=network-online.target
Wants=network-online.target

[Service]
User={{ vector_user }}
Group={{ vector_group }}
ExecStart={{ vector_binary_path }} --config {{ vector_config_dir }}/vector.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Environment
Environment=VECTOR_LOG_DIR={{ vector_log_dir }}
Environment=VECTOR_DATA_DIR={{ vector_data_dir }}

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths={{ vector_data_dir }} {{ vector_config_dir }} {{ vector_log_dir }}

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
